# misc/hades.target
[Unit]
Description=Boot into HadesOS
Requires=basic.target lxd.service ssh.service cron.service
AllowIsolate=yes
