# misc/hades-init.service
[Unit]
Description=hades-lock

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hades init

[Install]
WantedBy=multi-user.target
